# apps/photobooth/crontab
* * * * * find $HOME/public/pic -name *.jpg -mmin +3 -delete
